# crontab
0 1 * * * python3 /usr/src/environments/easyimovel/services/olx/main.py get >> /var/log/cron.log 2>&1

